# paper-restart.timer
[Unit]
Description=minecraft Server Restart Schedule Task

[Timer]
OnCalendar=*-*-* 04:59:30

[Install]
WantedBy=multi-user.target

