# paper-save.timer
[Unit]
Description=minecraft Server Save Schedule Task

[Timer]
OnBootSec=1min
OnUnitActiveSec=600s

[Install]
WantedBy=multi-user.target

